# etc/system/app.service
[Unit]
Description=Web App Service

[Service]
WorkingDirectory=/home/pi/livesectional
ExecStart=/usr/bin/uwsgi -i /home/pi/livesectional/uwsgi.ini
SyslogIdentifier=SRV_APPS

[Install]
WantedBy=multi-user.target